# systemd/memcontext.service
# Systemd 服务配置文件
# 用于在 Linux 服务器上运行 MemContext-Ad 作为系统服务

# 安装方法:
# 1. sudo cp systemd/memcontext.service /etc/systemd/system/
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable memcontext
# 4. sudo systemctl start memcontext

[Unit]
Description=MemContext-Ad Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/memcontext-ad

# 环境变量（或从文件加载）
EnvironmentFile=/opt/memcontext-ad/.env

# 启动命令
ExecStart=/opt/memcontext-ad/venv/bin/gunicorn \
    --bind 0.0.0.0:5019 \
    --workers 4 \
    --threads 2 \
    --timeout 300 \
    --access-logfile /var/log/memcontext/access.log \
    --error-logfile /var/log/memcontext/error.log \
    memdemo.app:app

# 重启策略
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65536
MemoryLimit=4G

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/memcontext-ad/memdemo/data /opt/memcontext-ad/files

[Install]
WantedBy=multi-user.target

